# packages/csr_azure_guestshell/csr_azure_guestshell-1.1.4.tar.gz/csr_azure_guestshell-1.1.4/csr_azure_guestshell/TokenMgr/auth-token.service
[Unit]
Description=Authentication Token service
After=waagent.service network.target

[Service]
Type=simple
User=guestshell
Group=guestshell
WorkingDirectory=/home/guestshell/azure/tools/TokenMgr
ExecStart=/home/guestshell/.local/lib/python2.7/site-packages/csr_azure_guestshell/TokenMgr/token_svr.py
Restart=always
RestartSec=1
StandardOutput=null
StandardError=syslog

[Install]
WantedBy=multi-user.target
