# packages/mixlight/mixlight-0.62-py3-none-any.whl/mixlight/hass.service
[Unit]
Description=HAss Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hass -c <folder>

[Install]
WantedBy=multi-user.target

