# packages/mixlight/mixlight-0.62-py3-none-any.whl/mixlight/mixlight.service
[Unit]
Description=MixLight Service 
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mixlight start <folder>

[Install]
WantedBy=multi-user.target
